# systemd/system/fstrim.timer
# /usr/lib/systemd/system/fstrim.timer
# Working path on Bunsenlabs:
# /etc/systemd/fstrim.timer


[Unit]
Description=Discard unused blocks at boot
Documentation=man:fstrim

[Timer]
OnCalendar=daily
#OnBootSec=1min
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
